# rootfs/etc/services.d/backend/run
#!/usr/bin/with-contenv bashio
export PORT=8099
export SUPERVISOR_TOKEN="${SUPERVISOR_TOKEN}"
export PYTHONPATH=/app/backend
exec python3 /app/backend/app.py
